# systemd/turnitin-celery.service
[Unit]
Description=Turnitin Bypass Celery Workers
After=network.target redis.service
Wants=redis.service

[Service]
Type=forking
User=YOUR_USERNAME
Group=YOUR_GROUP
WorkingDirectory=/workspaces/vision-computer
Environment="PATH=/usr/local/bin:/usr/bin:/bin"
Environment="PYTHONUNBUFFERED=1"

# Celery command
ExecStart=/usr/local/bin/celery -A app.celery_app worker \
    --loglevel=info \
    --concurrency=4 \
    --pool=prefork \
    --queues=unified,analysis,matching,bypass \
    --max-tasks-per-child=10 \
    --time-limit=600 \
    --soft-time-limit=540 \
    --pidfile=/workspaces/vision-computer/pids/celery_worker.pid \
    --logfile=/workspaces/vision-computer/logs/celery_worker.log \
    --detach

# PID file
PIDFile=/workspaces/vision-computer/pids/celery_worker.pid

# Stop command (graceful shutdown)
ExecStop=/bin/kill -TERM $MAINPID

# Restart policy
Restart=always
RestartSec=10
StartLimitInterval=200
StartLimitBurst=5

# Kill settings untuk celery
KillMode=mixed
KillSignal=SIGTERM
TimeoutStopSec=60

# Security settings
NoNewPrivileges=true
PrivateTmp=true

# Logging
StandardOutput=append:/workspaces/vision-computer/logs/celery.log
StandardError=append:/workspaces/vision-computer/logs/celery_error.log

[Install]
WantedBy=multi-user.target
