# systemd/turnitin-api.service
[Unit]
Description=Turnitin Bypass API (FastAPI + Gunicorn)
After=network.target redis.service
Wants=redis.service

[Service]
Type=notify
User=YOUR_USERNAME
Group=YOUR_GROUP
WorkingDirectory=/workspaces/vision-computer
Environment="PATH=/usr/local/bin:/usr/bin:/bin"
Environment="PYTHONUNBUFFERED=1"

# Gunicorn command
ExecStart=/usr/local/bin/gunicorn app.main:app \
    --workers 4 \
    --worker-class uvicorn.workers.UvicornWorker \
    --bind 0.0.0.0:8000 \
    --timeout 300 \
    --keep-alive 5 \
    --log-level info \
    --access-logfile /workspaces/vision-computer/logs/access.log \
    --error-logfile /workspaces/vision-computer/logs/error.log

# Restart policy
Restart=always
RestartSec=10
StartLimitInterval=200
StartLimitBurst=5

# Security settings
NoNewPrivileges=true
PrivateTmp=true

# Logging
StandardOutput=append:/workspaces/vision-computer/logs/api.log
StandardError=append:/workspaces/vision-computer/logs/api_error.log

[Install]
WantedBy=multi-user.target
